# ansible/roles/system_core/templates/backup.service.j2
[Unit]
Description=Daily System Backup
RequiresMountsFor={{ backup_path }}

[Service]
Type=oneshot
ExecStart=/usr/bin/tar -czf {{ backup_path }}/etc-backup.tar.gz /etc/fstab
