# ansible/roles/system_core/templates/backup.timer.j2
[Unit]
Description=Trigger Backup Job Daily

[Timer]
OnBootSec=5min
OnCalendar={{ backup_schedule }}
Persistent=true
Unit=backup-job.service

[Install]
WantedBy=timers.target
